# bootstrap/hadoop/supervisor/hadoop.conf
[program:hadoop]
command=/usr/bin/start-hdfs.sh
autostart=true
priority=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
