# bootstrap/hadoop/supervisor/hbase.conf
[program:hadoop]
command=/usr/bin/start-hbase.sh
autostart=true
priority=2
